# src/systemd/wg-quick.target
[Unit]
Description=WireGuard Tunnels via bwg-quick(8)
